# launchagents/com.user.openclaw-menubar.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.openclaw-menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>
        <string>__HOME__/bin/openclaw-menubar.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw-menubar.log</string>
    <!-- Required for GUI/menu bar apps -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
